# admin/уязвимость которая давно всем известна/docker-exim.service
[Unit]

Description=docker form vulner exim
After=docker.service
Requires=docker.service



[Service]
TimeoutStartSec=0
Restart=always
ExecStart=/usr/bin/docker run --rm --name exim -p 10.5.11.12:25:25 --cap-add=SYS_PTRACE --security-opt seccomp=unconfined exim:suid
ExecStop=/usr/bin/docker kill exim

[Install]

WantedBy=default.target

